# full/host/systemd/vps-mcp-reboot.timer
[Unit]
Description=Weekly automatic reboot timer for the VPS MCP host
Documentation=https://github.com/paijp/vps-mcp

[Timer]
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=600
Persistent=false

[Install]
WantedBy=timers.target
